# installer/linux/debian/prerm
#!/bin/bash
set -e

# 卸载前停止并移除服务
if systemctl is-active --quiet lark-agent-bridge 2>/dev/null; then
    systemctl stop lark-agent-bridge || true
fi
if systemctl is-enabled --quiet lark-agent-bridge 2>/dev/null; then
    systemctl disable lark-agent-bridge || true
fi
rm -f /etc/systemd/system/lark-agent-bridge.service
systemctl daemon-reload 2>/dev/null || true
